# srcs/influxdb/supervisord.conf
[supervisord]
nodaemon = true
[program:influxdb]
command =  sh -c "influxd -config /etc/influxdb.conf && kill $(cat supervisord.pid)"


[program:telegraf]
command =  sh -c "/telegraf/usr/bin/telegraf && kill $(cat supervisord.pid)"

