# srcs/wordpress/supervisord.conf
[supervisord]
nodaemon = true
#супервизор не любит когда запускают демона ( php-fpm7 - демон) с рутом без указания юзера в конфиге, но можно запустить фпм как процесс, для этого надо добавить --nodaemonize
[program:php-fpm]
command = sh -c "php-fpm7 --nodaemonize && kill $(cat supervisord.pid)"

[program:nginx]
command = sh -c "nginx  -g 'daemon off;' && kill $(cat supervisord.pid)"